# casa.service
[Unit]
Description=Casa smart home

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/casa

[Install]
WantedBy=multi-user.target
